# roboticscape_service/roboticscape.service
[Unit]
Description=roboticscape
Requires=systemd-modules-load.service

[Service]
User=root
PIDFile=/var/run/robotcscape.pid
ExecStartPre=/usr/bin/rc_startup_routine
ExecStart=/etc/roboticscape/link_to_startup_program
ExecStop=/usr/bin/rc_kill

[Install]
WantedBy=multi-user.target
