# debian/prerm
#!/bin/sh


# stop and disable the battery monitor service
systemctl stop rc_battery_monitor
systemctl disable rc_battery_monitor

# stop and disable the robotcontrol service
systemctl stop robotcontrol
systemctl disable robotcontrol

#DEBHELPER#

exit 0
